# config/crontab
@hourly bundle exec rails runner 'DlCleanup.call'
